# hades.conf
[program:hades_worker]
command = python3 worker.py -ip 139.219.138.34
directory = /home/hades/hades
user = root
autostart = true
autorestart = true