# hades.conf
[program:hades]
command = python3.7 server.py 0.0.0.0 8102 139.219.138.34:8102
directory = /root/hades/hades
user = root
autostart = true
autorestart = true
